# systemd/fmm.service
[Unit]
Description=A simple Atom and RSS feed to email 'converter'

[Service]
Type=simple
ExecStart=/usr/bin/fmm
User=fmm
